# nxlog-ce-3/packaging/debian/prerm
#!/bin/bash

set -e

case "$1" in

	remove)
	    if systemctl is-active --quiet nxlog; then
		systemctl stop nxlog >/dev/null 2>&1
		systemctl disable nxlog
	    fi
	;;

	*)
	;;

esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0

